# templates/read_currentcost.conf
#!upstart

exec /usr/local/electricity_logging/read_currentcost.py yes

start on runlevel [2345]
stop on runlevel [016]

respawn

